# sysfiles/profile.d/20-shinit.sh
if [ -r /etc/shinit/env.sh ]; then
   . /etc/shinit/env.sh || :
fi
